# systemd/notifyd.service
[Unit]
Description=Notify Notification Service
After=network.target

[Service]
ExecStart=/usr/bin/notifyd --config /etc/notify/config.yaml
Restart=always
User=root
Group=notify

[Install]
WantedBy=multi-user.target
